# scripts/systemd/omnimind-backend-direct.service
[Unit]
Description=OmniMind Backend API (Direct)
After=network.target omnimind-qdrant.service
Wants=omnimind-qdrant.service

[Service]
Type=simple
User=fahbrain
Group=fahbrain
WorkingDirectory=/home/fahbrain/projects/omnimind
Environment="PYTHONUNBUFFERED=1"
Environment="OMNIMIND_WORKSPACE=/home/fahbrain/projects/omnimind"
Environment="OMNIMIND_LOG_LEVEL=INFO"

# Virtual environment activation
ExecStartPre=/bin/bash -c 'source /home/fahbrain/projects/omnimind/.venv/bin/activate'

# Main service
ExecStart=/home/fahbrain/projects/omnimind/.venv/bin/python -m uvicorn web.backend.main:app --host 127.0.0.1 --port 8000 --log-level info

# Restart policy
Restart=on-failure
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/fahbrain/projects/omnimind/logs /home/fahbrain/projects/omnimind/data /home/fahbrain/projects/omnimind/.omnimind

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=omnimind-backend

[Install]
WantedBy=multi-user.target

