# scripts/systemd/omnimind-frontend.service
[Unit]
Description=OmniMind Frontend Dashboard
After=network.target omnimind.service
Requires=omnimind.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/fahbrain/projects/omnimind
EnvironmentFile=/home/fahbrain/projects/omnimind/.env
# Iniciar apenas o frontend, sem dependências (backend já está rodando via systemd)
ExecStart=/usr/bin/docker-compose -f deploy/docker-compose.yml up --no-deps frontend
ExecStop=/usr/bin/docker-compose -f deploy/docker-compose.yml stop frontend
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
